# mtpgo_openrc
#!/sbin/openrc-run

name="mtpgo"
description="Mtproxygo Service"
command="/bin/sh"
command_args="-c '/var/mtpgo/mtpgo -c /var/mtpgo/config.ini'"
command_background="yes"
pidfile="/var/run/${RC_SVCNAME}.pid"
start_stop_daemon_args="--user root:root"
